# board/sdn/h1502/rootfs-overlay/etc/profile.d/console.sh
#!/bin/sh
export PS1='$(whoami)@${HOSTNAME%%.*}:$(pwd)# '

